# files/system/etc/systemd/system/systemd-remount-fs.service.d/override.conf
[Unit]
# Skip remount attempts on immutable/OSTree systems where / is RO
ConditionPathIsReadWrite=/
