# files/system/etc/systemd/system/ext.automount
[Unit]
Description=Automount for /ext

[Automount]
Where=/ext

[Install]
WantedBy=multi-user.target
