# files/system/etc/systemd/system/proxmox-backup.timer
[Unit]
Description=Nightly PBS backup

[Timer]
OnCalendar=*-*-* 01:50:00
Persistent=true
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
