# files/system/etc/systemd/system/storage-health.service
[Unit]
Description=Check ZFS + MD RAID + SMART/NVMe, send email and log to syslog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=HOME=/root
ExecStart=/usr/libexec/storage-health-check.sh
