# files/system/etc/systemd/system/swapfile-setup.service
# /etc/systemd/system/swapfile-setup.service
[Unit]
Description=Prepare /var/swapfile for disk-based swap (conditional; Btrfs-safe)
DefaultDependencies=no
RequiresMountsFor=/var
Before=var-swapfile.swap
After=local-fs-pre.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/swapfile-setup
RemainAfterExit=yes

[Install]
WantedBy=var-swapfile.swap
