# files/system/etc/systemd/system/swapfile-setup.service
[Unit]
Description=Prepare /var/swapfile for disk-based swap (conditional; Btrfs-safe)
After=local-fs.target
RequiresMountsFor=/var
Before=var-swapfile.swap

[Service]
Type=oneshot
ExecStart=/usr/libexec/swapfile-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
