# files/system/etc/systemd/system/work.automount
[Unit]
Description=Automount for /work

[Automount]
Where=/work

[Install]
WantedBy=multi-user.target
