# files/system/lib/systemd/system/work.mount
[Mount]
What=/home/dmoseley/work
Where=/work
Type=none
Options=bind
